# systemd/flatpak-upgrade.service
[Unit]
Description=flatpak-upgrade
[Service]
User=empthollow
ExecStart=flatpak update -y
