# systemd/pacman-upgrade.service
[Unit]
Description=pacman-upgrade
[Service]
User=root
ExecStart=pacman -Syu --noconfirm
ExecStartPost=paccache -ruk1
