# systemd/pacman-sync.service
[Unit]
Description=pacman-sync
[Service]
User=root
ExecStart=pacman -Syy
