# systemd/spotify-restart.service
[Unit]
Description=spotify-restart
[Service]
User=empthollow
ExecStart=killall spotify && sleep 2 && flatpak run com.spotify.Client
